# metroclock.service
[Unit]
Description=MetroClock LED Display
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/zaid/metroclock
ExecStart=/usr/bin/python3 /home/zaid/metroclock/main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=METROCLOCK_CONFIG_PATH=/etc/metroclock/config.json
EnvironmentFile=-/etc/metroclock/secrets.env

[Install]
WantedBy=multi-user.target
